# deploy/systemd/satellite-environment-downloader.service
[Unit]
Description=Satellite Environment Downloader (Streamlit)
After=network.target

[Service]
Type=simple
User=CHANGE_ME
WorkingDirectory=CHANGE_ME/satellite_dengue_excel
ExecStart=CHANGE_ME/satellite_dengue_excel/scripts/start_streamlit.sh
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
